# target/common/rootfs/etc/profile.d/helloworld.sh
printf "\n\e[4;37m%s\e[1m\e[38;5;208m%s\e[38;5;209m%s\e[38;5;210m%s\e[38;5;211m%s\e[38;5;212m%s\e[38;5;177m%s\e[38;5;141m%s\e[38;5;105m%s\e[m\n\n" "[PULP Platform] Welcome to Linux on " "C" "h" "e" "s" "h" "i" "r" "e"